# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/tmp/supervisord.pid

[program:postgres]
command=/usr/lib/postgresql/14/bin/postgres -D /var/lib/postgresql/data -c listen_addresses=127.0.0.1
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/postgres.log
stderr_logfile=/var/log/supervisor/postgres_err.log
priority=1
startsecs=5

[program:backend]
command=node /app/backend/src/server.js
directory=/app/backend
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/backend.log
stderr_logfile=/var/log/supervisor/backend_err.log
environment=NODE_ENV="production",PORT="5000",DB_HOST="127.0.0.1",DB_PORT="5432",DB_NAME="taskflow",DB_USER="taskflow_user",DB_PASSWORD="taskflow_hf_secret",JWT_SECRET="taskflow_huggingface_jwt_secret_change_for_production_use",JWT_EXPIRES_IN="7d",FRONTEND_URL="http://localhost:7860"
priority=20
startsecs=3

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx_err.log
priority=30
startsecs=2
